# autochown.service
[Unit]
Description=Autochown

[Service]
ExecStart=/usr/bin/autochown /etc/autochownd.conf

[Install]
WantedBy=multi-user.target
